# templates/kibana.legacy.service.lsbheader.erb
#!/bin/sh
#
# kibana    Web interface to Logstash.
#
# chkconfig:   345 99 01
# description: Starts and stops a single kibana instance on this system.

### BEGIN INIT INFO
# Provides: kibana
# Required-Start: $all
# Required-Stop: $all
# Default-Start: 2 3 4 5
# Default-Stop: 0 1 6
# Short-Description: This service manages the kibana daemon
# Description: Kibana is the Elasticsearch WebUI for time-based data
### END INIT INFO

DESC="Kibana"
NAME=kibana

